# 01_etl/crontab
*/1 * * * * /usr/bin/python3.8 /root/main.py >> /var/log/cron.log 2>&1
